# infra/potager-prod-bot.service
[Unit]
Description=Assistant Potager — Bot Telegram (PROD)
After=network.target postgresql.service potager-prod.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/potager-prod
EnvironmentFile=/opt/potager-prod/.env.prod
Environment=APP_ENV=prod
ExecStart=/opt/potager-prod/venv/bin/python bot.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
