# infra/potager-bot.service
[Unit]
Description=Assistant Potager — Bot Telegram
After=network.target postgresql.service potager.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/potager
EnvironmentFile=/opt/potager/.env.prod
Environment=APP_ENV=prod
ExecStart=/opt/potager/venv/bin/python bot.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
